# etc/supervisor-labsome.conf
[supervisord]
nodaemon=true

[program:labsome-server]
command=/usr/local/bin/labsome-server
autorestart=true
stdout_logfile=/var/log/labsome.log
redirect_stderr=true
