# etc/supervisor-labsome.conf
[supervisord]
nodaemon = true
logfile = /var/log/supervisord.log
pidfile = /var/run/supervisord.pid

[program:api]
command = uwsgi --socket 127.0.0.1:5000 -w "labsome.backend.app:create_app()" --python-auto-reload 1
autorestart = true
user = www-data
stdout_logfile = /var/log/labsome/api.log
redirect_stderr = true

[program:nginx]
command = nginx -c /opt/labsome/etc/nginx.conf
autorestart = true
stdout_logfile = /var/log/labsome/nginx.log
redirect_stderr = true
