# dist/etc/supervisor-labsome.conf
[supervisord]
nodaemon = true
logfile = /var/log/supervisord.log
pidfile = /var/run/supervisord.pid

[program:api-server]
command = gunicorn --bind 127.0.0.1:5000 --worker-class eventlet --access-logfile - "labsome_api.app:create_app_with_console_logging()"
autorestart = true
user = www-data
stdout_logfile = /var/log/labsome/api.log
redirect_stderr = true

[program:notify-server]
command = labsome-notify-server
autorestart = true
user = www-data
stdout_logfile = /var/log/labsome/notifications.log
redirect_stderr = true

[program:nginx]
command = nginx -c /opt/labsome/dist/etc/nginx.conf
autorestart = true
stdout_logfile = /var/log/labsome/nginx.log
redirect_stderr = true
